# daemon/ai.muavin.cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.muavin.cron</string>
    <key>ProgramArguments</key>
    <array>
        <string>__MUAVIN_BIN__/muavin-cron</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HOME__/.muavin</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:__HOME__/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/muavin-cron.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/muavin-cron.error.log</string>
</dict>
</plist>
